# templates/etc/systemd/system/docker.service.d/service-overrides.conf.erb
[Service]
EnvironmentFile=-/etc/sysconfig/docker
EnvironmentFile=-/etc/sysconfig/docker-storage
ExecStart=/usr/bin/dockerd $OPTIONS \
      $DOCKER_STORAGE_OPTIONS
