# templates/etc/systemd/system/docker.service.d/service-overrides.conf.erb
[Service]
<% if @service_config -%>
EnvironmentFile=-<%= @service_config %>
<% end -%>
<% if @storage_config -%>
EnvironmentFile=-<%= @storage_config %>
<% end -%>
ExecStart=/usr/bin/dockerd $OPTIONS \
      $DOCKER_STORAGE_OPTIONS
